# rootfs/etc/services.d/naim-bridge/run
#!/usr/bin/with-contenv bashio

export QUTE_IP=$(bashio::config 'qute_ip')
export VOLUME_CINEMA=$(bashio::config 'volume_cinema')
export VOLUME_SPOTIFY=$(bashio::config 'volume_spotify')
export SPOTIFY_CLIENT_ID=$(bashio::config 'spotify_client_id')
export SPOTIFY_CLIENT_SECRET=$(bashio::config 'spotify_client_secret')
export SPOTIFY_DEVICE_NAME=$(bashio::config 'spotify_device_name')

bashio::log.info "Demarrage Naim Bridge sur ${QUTE_IP}"

cd /usr/bin
bashio::log.info "Lancement Python..."
exec python3 -u bridge.py
