# rootfs/etc/services.d/naim-bridge/run
#!/command/execlineb -P
cd /usr/bin
python3 -u bridge.py
