# rootfs/etc/services.d/naim-bridge/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Test service démarré"
exec sleep infinity

